# InitSystems/daemon/systemd/myservice.service
[Service]
ExecStart=/home/darklord/Code/daemon/systemd/myservice_script.sh


